# roles/rpi_mxbi_baseline/templates/mediamtx.service.j2
# Managed by Ansible. Do not edit this file on the target.
[Unit]
Description=MediaMTX real-time media server
Wants=network-online.target
After=network-online.target

[Service]
ExecStart={{ rpi_mxbi_baseline_homebrew_prefix }}/bin/mediamtx %h/.config/mediamtx/mediamtx.yml
WorkingDirectory=%h
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
